# userspace/var/lib/pacman/local/libteam-1.31-3/files
%FILES%
usr/
usr/bin/
usr/bin/bond2team
usr/bin/teamd
usr/bin/teamdctl
usr/bin/teamnl
usr/include/
usr/include/team.h
usr/include/teamdctl.h
usr/lib/
usr/lib/libteam.so
usr/lib/libteam.so.5
usr/lib/libteam.so.5.6.1
usr/lib/libteamdctl.so
usr/lib/libteamdctl.so.0
usr/lib/libteamdctl.so.0.1.5
usr/lib/pkgconfig/
usr/lib/pkgconfig/libteam.pc
usr/lib/pkgconfig/libteamdctl.pc
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/team-1.0-py3.9.egg-info
usr/lib/python3.9/site-packages/team/
usr/lib/python3.9/site-packages/team/__init__.py
usr/lib/python3.9/site-packages/team/__pycache__/
usr/lib/python3.9/site-packages/team/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/team/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/team/__pycache__/core.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/team/__pycache__/core.cpython-39.pyc
usr/lib/python3.9/site-packages/team/_capi.cpython-39-aarch64-linux-gnu.so
usr/lib/python3.9/site-packages/team/core.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/teamd@.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/teamd.conf
usr/share/doc/
usr/share/doc/libteam/
usr/share/doc/libteam/api/
usr/share/doc/libteam/api/bc_s.png
usr/share/doc/libteam/api/bdwn.png
usr/share/doc/libteam/api/closed.png
usr/share/doc/libteam/api/deprecated.html
usr/share/doc/libteam/api/dir_3108ee1453c90726cc7ff8c1a8003399.html
usr/share/doc/libteam/api/dir_b684ca15ef9da6686f0dc3b9dd11b939.html
usr/share/doc/libteam/api/doc.png
usr/share/doc/libteam/api/doxygen.css
usr/share/doc/libteam/api/doxygen.svg
usr/share/doc/libteam/api/dynsections.js
usr/share/doc/libteam/api/folderclosed.png
usr/share/doc/libteam/api/folderopen.png
usr/share/doc/libteam/api/group__core.html
usr/share/doc/libteam/api/group__ifinfo.html
usr/share/doc/libteam/api/group__libteam.html
usr/share/doc/libteam/api/group__libteamdctl.html
usr/share/doc/libteam/api/group__option.html
usr/share/doc/libteam/api/group__ports.html
usr/share/doc/libteam/api/group__stringify.html
usr/share/doc/libteam/api/index.html
usr/share/doc/libteam/api/jquery.js
usr/share/doc/libteam/api/menu.js
usr/share/doc/libteam/api/menudata.js
usr/share/doc/libteam/api/modules.html
usr/share/doc/libteam/api/nav_f.png
usr/share/doc/libteam/api/nav_g.png
usr/share/doc/libteam/api/nav_h.png
usr/share/doc/libteam/api/open.png
usr/share/doc/libteam/api/pages.html
usr/share/doc/libteam/api/search/
usr/share/doc/libteam/api/search/all_0.html
usr/share/doc/libteam/api/search/all_0.js
usr/share/doc/libteam/api/search/all_1.html
usr/share/doc/libteam/api/search/all_1.js
usr/share/doc/libteam/api/search/all_2.html
usr/share/doc/libteam/api/search/all_2.js
usr/share/doc/libteam/api/search/all_3.html
usr/share/doc/libteam/api/search/all_3.js
usr/share/doc/libteam/api/search/all_4.html
usr/share/doc/libteam/api/search/all_4.js
usr/share/doc/libteam/api/search/close.svg
usr/share/doc/libteam/api/search/functions_0.html
usr/share/doc/libteam/api/search/functions_0.js
usr/share/doc/libteam/api/search/groups_0.html
usr/share/doc/libteam/api/search/groups_0.js
usr/share/doc/libteam/api/search/groups_1.html
usr/share/doc/libteam/api/search/groups_1.js
usr/share/doc/libteam/api/search/groups_2.html
usr/share/doc/libteam/api/search/groups_2.js
usr/share/doc/libteam/api/search/groups_3.html
usr/share/doc/libteam/api/search/groups_3.js
usr/share/doc/libteam/api/search/mag_sel.svg
usr/share/doc/libteam/api/search/nomatches.html
usr/share/doc/libteam/api/search/pages_0.html
usr/share/doc/libteam/api/search/pages_0.js
usr/share/doc/libteam/api/search/search.css
usr/share/doc/libteam/api/search/search.js
usr/share/doc/libteam/api/search/search_l.png
usr/share/doc/libteam/api/search/search_m.png
usr/share/doc/libteam/api/search/search_r.png
usr/share/doc/libteam/api/search/searchdata.js
usr/share/doc/libteam/api/splitbar.png
usr/share/doc/libteam/api/sync_off.png
usr/share/doc/libteam/api/sync_on.png
usr/share/doc/libteam/api/tab_a.png
usr/share/doc/libteam/api/tab_b.png
usr/share/doc/libteam/api/tab_h.png
usr/share/doc/libteam/api/tab_s.png
usr/share/doc/libteam/api/tabs.css
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/bond2team.1.gz
usr/share/man/man5/
usr/share/man/man5/teamd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/teamd.8.gz
usr/share/man/man8/teamdctl.8.gz
usr/share/man/man8/teamnl.8.gz

